# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0777
chown=root:root

[supervisord]
logfile=/var/log/supervisor/supervisor.log
logfile_maxbytes=50MB
loglevel=error
pidfile=/var/run/supervisor.pid
nodaemon=false
user=root
childlogdir=/var/log/supervisor/

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:hass]
command=/usr/bin/hactl run
autostart=false
autorestart=false
user=root
